# com.guoxl.war.timeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.guoxl.war.timeline</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>/Users/guoxl/.openclaw/war_timeline/run_update.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/guoxl/.openclaw/war_timeline</string>
  <key>StartInterval</key>
  <integer>43200</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/guoxl/.openclaw/war_timeline/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/guoxl/.openclaw/war_timeline/logs/launchd.err.log</string>
</dict>
</plist>
